# systemd/wal-cache-apply.service
[Unit]
Description=Apply pywal when wallpaper cache changes
ConditionUser=!root

[Service]
Type=oneshot
ExecStart=%h/.local/bin/wal-apply-cache
